# archiso/airootfs/etc/systemd/system/vibe-installer.service
[Unit]
Description=Vibe-Linux Guided Installer
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/local/bin/vibe-installer-bootstrap

[Service]
Type=simple
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=tty
StandardError=tty
ExecStart=/usr/local/bin/vibe-installer-bootstrap
Restart=no

[Install]
WantedBy=multi-user.target

